# deploy/systemd/slither-agent.service
[Unit]
Description=Slither EDR agent (eBPF-based)
Documentation=https://github.com/t3rmit3/slither
After=network-online.target
Wants=network-online.target
ConditionPathExists=/sys/kernel/btf/vmlinux

[Service]
Type=simple
ExecStart=/usr/local/bin/slither-agent --config /etc/slither/agent.yaml
ExecReload=/bin/kill -HUP $MAINPID

# Run as root. The agent needs to load BPF programs, read every process's
# /proc entry, and hash arbitrary files for executable provenance. A
# dedicated unprivileged user buys nothing when all four capabilities below
# are effectively root-equivalent, and Phase 1 sticks to root + capability
# bounding per IMPLEMENTATION.md §3.8.
User=root
Group=root

# Cap the kernel privileges the agent can ever acquire. The kernel drops
# anything outside this set on first exec, so a compromised binary cannot
# regain, e.g., CAP_NET_ADMIN or CAP_SYS_MODULE.
#   CAP_BPF            — load/attach BPF programs and maps (>=5.8).
#   CAP_PERFMON        — required alongside CAP_BPF for tracepoint/kprobe
#                        attach on 5.8+ kernels.
#   CAP_SYS_PTRACE     — read /proc/<pid>/exe + cmdline for enrichment
#                        across user namespaces.
#   CAP_DAC_READ_SEARCH — open executables for hashing regardless of DAC
#                        permissions (hashing worker pool).
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_PTRACE CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_PTRACE CAP_DAC_READ_SEARCH

# NoNewPrivileges stays off as a belt-and-braces default. Historically some
# 5.x kernels refused bpf(BPF_PROG_LOAD) from a task with no_new_privs set
# when the program type required CAP_PERFMON (the LSM hook re-evaluated
# credentials against the current restricted set rather than the bounding
# set). Our 5.15 floor and the supported matrix (Debian 13, RHEL 10,
# Ubuntu 22.04+) no longer exhibit this; flipping to yes is deferred until
# CI proves it across the matrix.
NoNewPrivileges=no

# Filesystem hygiene. We still want the agent contained even without
# NoNewPrivileges — these hardenings are independent of the capability
# check above.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/slither /var/log/slither
RuntimeDirectory=slither
StateDirectory=slither
LogsDirectory=slither
PrivateTmp=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes

# Restart policy: the agent is expected to exit on fatal config/rule errors
# (see §3.5). Let systemd back off and retry so an operator can fix the
# config and `systemctl reload` or let the restart pick it up.
Restart=on-failure
RestartSec=5s

# Resource ceilings. BPF ringbufs + hashing worker pool sit comfortably
# under these on a 50–500-host fleet.
LimitNOFILE=65536
LimitMEMLOCK=infinity
TasksMax=1024

[Install]
WantedBy=multi-user.target
